# desktop-database/.config/systemd/user/desktop-database.path
[Unit]
Description=Watch for changes in $XDG_DATA_HOME/applications/

[Path]
MakeDirectory=true
PathChanged=%h/.local/share/applications

[Install]
WantedBy=paths.target
