# rust/.config/systemd/user/cargo-bash-completion.path
[Unit]
Description=Watch for changes in ~/.opt/cargo/bin/cargo

[Path]
PathChanged=%h/.opt/cargo/bin/cargo

[Install]
WantedBy=paths.target
