# man/.config/systemd/user/mandb.path
[Unit]
Description=Watch for changes in $XDG_DATA_HOME/man

[Path]
MakeDirectory=true
# %D is not available on older versions of systemd, like on Amazon Linux
# PathChanged=%D/man
PathChanged=%h/.local/share/man

[Install]
WantedBy=paths.target
