# man/.config/systemd/user/mandb.path
[Unit]
Description=Watch for changes in $XDG_DATA_HOME/man

[Path]
MakeDirectory=true
PathChanged=%D/man

[Install]
WantedBy=paths.target
