# aerc/.config/systemd/user/aerc-accounts-conf.path
[Unit]
Description=Watch for changes in the aerc accounts.conf file

[Path]
PathChanged=%E/aerc/accounts.conf

[Install]
WantedBy=paths.target
