# flyctl/.config/systemd/user/flyctl-bash-completion.path
[Unit]
Description=Watch for changes in ~/.opt/fly/bin/flyctl

[Path]
PathChanged=%h/.opt/fly/bin/flyctl

[Install]
WantedBy=paths.target
