# gh/.config/systemd/user/gh-bash-completion.path
[Unit]
Description=Watch for changes in gh

[Path]
PathChanged=/usr/bin/gh

[Install]
WantedBy=paths.target
