# systemd/.config/systemd/user/mandb.path
[Unit]
Description=Watch for changes in XDG_DATA_HOME/man

[Path]
MakeDirectory=true
# systemd, why can't I use environment variables!!!!
PathChanged=%h/.local/share/man

[Install]
WantedBy=paths.target
