# systemd/.config/systemd/user/mandb.timer
[Unit]
Description=Run mandb on a schedule

[Timer]
OnUnitActiveSec=1day

[Install]
WantedBy=timers.target
