# rustup/.config/systemd/user/rustup-bash-completion.path
[Unit]
Description=Watch for changes in ~/.opt/cargo/bin/rustup

[Path]
PathChanged=%h/.opt/cargo/bin/rustup

[Install]
WantedBy=paths.target
