# deno/.config/systemd/user/deno-bash-completion.path
[Unit]
Description=Watch for changes in ~/.opt/deno/bin/deno

[Path]
PathChanged=%h/.opt/deno/bin/deno

[Install]
WantedBy=paths.target
