# uv/.config/systemd/user/uv-bash-completion.path
[Unit]
Description=Watch for changes in ~/.opt/cargo/bin/uv

[Path]
PathChanged=%h/.opt/cargo/bin/uv

[Install]
WantedBy=paths.target
