# uv/.config/systemd/user/uvx-bash-completion.path
[Unit]
Description=Watch for changes in ~/.opt/cargo/bin/uvx

[Path]
PathChanged=%h/.opt/cargo/bin/uvx

[Install]
WantedBy=paths.target
